# roles/TG/templates/supervisor/conf.d/weedfs.conf
[program:weedfs_server]
command=/bin/bash ./run_server.sh /home/{{ansible_user}}/platformData/weedfs/master/ 9333
process_name=%(program_name)s
numprocs=1
directory=/home/{{ansible_user}}/platformTG/weedfs/latest
;umask=022
priority=700
autostart=true
autorestart=true
;startsecs=3
startretries=3
exitcodes=0,2
;stopsignal=TERM
stopwaitsecs=5
stopasgroup=true
killasgroup=true
;user=root
redirect_stderr=true
stdout_logfile=/home/{{ansible_user}}/platformData/logs/%(program_name)s.out
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

[program:weedfs_volume]
command=/bin/bash ./run_volume.sh /home/{{ansible_user}}/platformData/weedfs/volume%(process_num)02d localhost:9333 93%(process_num)02d
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/home/{{ansible_user}}/platformTG/weedfs/latest
;umask=022
priority=888
autostart=true
autorestart=true
;startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
exitcodes=0,2
;stopsignal=TERM
stopwaitsecs=10
;user=root
redirect_stderr=true
stdout_logfile=/home/{{ansible_user}}/platformData/logs/%(program_name)s_%(process_num)02d.out
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

[group:seaweedfs]
programs=weedfs_server,weedfs_volume
priority=999
